# chromeos-base/auto-expand-partition/files/auto-expand-partition.conf
# Copyright (c) 2017 Flint Innovations Limited. All rights reserved.
# Distributed under the terms of the GNU General Public License v2

description     "Auto expand stateful partition on first boot"
author          "kraml@flintos.io"

start on starting boot-services

task

script
  if [ ! -f /mnt/stateful_partition/.autoexpanded ]; then
    stateful_partition=$(findmnt -n -o source /mnt/stateful_partition)
    root_disk="/dev/$(lsblk -n -o pkname $stateful_partition)"

    echo "Expanding stateful partition..."
    start_sec=$(/usr/sbin/partx -s -b -g -o start $stateful_partition)
    last_byte=$(/bin/lsblk -b -d -n -o size $root_disk)
    last_sec=$(($last_byte / 512))

    # After write image to SD card, the 2nd GPT table/header are at the middle
    # of the disk, as the image is usually smaller than the capacity of the card.
    # The ChromeOS cgpt repair command can move them to the end of the disk but
    # later when running cgpt add command it won't allow expanding the partition
    # size beyond the last sector of the original image.
    # A possible guess is that cgpt repair command doesn't update the LBA# of
    # the 2nd table/header recorded in the primary header. This cause cgpt add
    # command (mistakenly) to think the expansion would overwrite the 2nd GPT
    # table/header (that looks like were still in the middle of the disk).
    # Thus sgdisk is used instead as it does the job perfectly.
    #
    /usr/sbin/sgdisk -e $root_disk

    # 2nd GPT table/header are at the end of the disk and occupy 32+1 sectors.
    # So the new partition can't take all spare space to the end of the disk.
    # Minus 33 sectors should be enough but a bit more won't hurt.
    #
    /usr/bin/cgpt add -i 1 -s $(($last_sec - $start_sec - 10240)) $root_disk

    # Tell the kernel to pick up the new partition size
    /usr/sbin/partx -u $stateful_partition

    # Enlarge filesystem
    /sbin/resize2fs $stateful_partition

    touch /mnt/stateful_partition/.autoexpanded
  fi
end script
